# service_files/host_power_state_check.service
[Unit]
Description=Host Power Check

[Service]
ExecStart=/usr/bin/env host_power_state_check.sh
Restart=no

[Install]
WantedBy=multi-user.target